# vps-agent/agent.service
[Unit]
Description=AI Trading Agent (Multi-Timeframe, Claude-powered)
Documentation=https://github.com/Siddikuk/Trading-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# SECURITY: run as a dedicated low-privilege user, NOT root.
# Before first deploy: useradd -r -s /bin/false -d /home/trading trading
#                      chown -R trading:trading /home/user/Trading-agent/vps-agent
User=trading
WorkingDirectory=/home/user/Trading-agent/vps-agent
EnvironmentFile=/home/user/Trading-agent/vps-agent/.env
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=30
# Give the agent up to 60s to finish current cycle on shutdown
TimeoutStopSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-agent

[Install]
WantedBy=multi-user.target
